# systemd/cleanup.timer
[Unit]
Description=Run sandbox cleanup daily

[Timer]
OnBootSec=5min
OnUnitActiveSec=1d
Unit=cleanup.service
Persistent=true

[Install]
WantedBy=timers.target
