# systemd/menace.service
[Unit]
Description=Menace Service Supervisor
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/menace
ExecStart=/usr/bin/python -m menace.service_supervisor
Environment="MENACE_BOOTSTRAP_WAIT_SECS=360"
Environment="MENACE_BOOTSTRAP_VECTOR_WAIT_SECS=540"
Environment="BOOTSTRAP_STEP_TIMEOUT=360"
Environment="BOOTSTRAP_VECTOR_STEP_TIMEOUT=540"
Environment="MENACE_BOOTSTRAP_STAGGER_SECS=30"
Environment="MENACE_BOOTSTRAP_STAGGER_JITTER_SECS=30"
Environment="MENACE_BOOTSTRAP_READINESS_MODE=staged"
ExecStartPre=/usr/bin/python -m bootstrap_conflict_check
Restart=always

[Install]
WantedBy=multi-user.target
