# systemd/menace.service
[Unit]
Description=Menace Service Supervisor
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/menace
ExecStart=/usr/bin/python -m menace.service_supervisor
Environment="MENACE_BOOTSTRAP_WAIT_SECS=120"
Environment="BOOTSTRAP_STEP_TIMEOUT=90"
Environment="BOOTSTRAP_VECTOR_STEP_TIMEOUT=150"
Restart=always

[Install]
WantedBy=multi-user.target
