# systemd/menace.service
[Unit]
Description=Menace Service Supervisor
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/menace
ExecStart=/usr/bin/python -m menace.service_supervisor
Restart=always

[Install]
WantedBy=multi-user.target
