# systemd/sandbox_purge.service
[Unit]
Description=Sandbox Runner stale resource cleanup
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/menace_sandbox
ExecStart=/usr/bin/python -m sandbox_runner.cli --purge-stale

[Install]
WantedBy=multi-user.target
